# systemd/offlineimap@.timer
[Unit]
Description=Crontjob for offlineimap

[Timer]
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target
